# image/dindnet.service
[Unit]
Description=DIND network setup

[Service]
ExecStart=/usr/local/bin/dindnet

[Install]
WantedBy=docker.service
